# deploy/systemd/simplepool-dashboard.service
[Unit]
Description=simplepool dashboard (read-only web UI)
Documentation=https://github.com/rsantacroce/simplepool
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@ROOT@/dashboard
Environment=PORT=8081
Environment=PROXY_DB_PATH=@ROOT@/data/shares.db
# --- PPS admin panel — /admin route ---
# When ADMIN_USER + ADMIN_PASSWORD are both set, /admin serves the pool
# operator view (reserve balance, ledger totals, per-worker owed,
# in-flight payouts). Without both, /admin returns 503 (fail-closed).
Environment=THUNDER_RPC_URL=http://127.0.0.1:6009
# Environment=POOL_THUNDER_RESERVE_ADDRESS=<paste from proxy.conf>
# Environment=ADMIN_USER=admin
# Environment=ADMIN_PASSWORD=<generate with `openssl rand -base64 21`>
# Must equal proxy.conf's pps_sats_per_diff — used by the per-worker
# audit view to cross-check the C proxy's accrual math.
Environment=POOL_PPS_SATS_PER_DIFF=1000
ExecStart=/usr/bin/node @ROOT@/dashboard/server.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadOnlyPaths=@ROOT@

[Install]
WantedBy=multi-user.target
